# systemd/sentinel-node-adapter.service
[Unit]
Description=SENTINEL Node Adapter — ESP32 to Sentinel topic bridge
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=<user>
WorkingDirectory=/home/<user>/Presence
EnvironmentFile=/home/<user>/Presence/.env
ExecStart=/usr/bin/python3 -m sentinel.adapters \
    --mqtt-host ${MQTT_HOST} --log-level ${LOG_LEVEL}
Restart=always
RestartSec=5
MemoryMax=128M
CPUQuota=15%

# Hardening
ProtectSystem=strict
ReadWritePaths=/home/<user>/Presence /home/<user>/sentinel
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
